# pi/nowify-watchdog.service
[Unit]
Description=Nowify kiosk watchdog (restart Chromium when stuck, never on expired token)
After=network-online.target

[Service]
Type=oneshot
# Set the relaunch command for YOUR setup (see pi/README.md). Examples:
#   Environment=NOWIFY_RESTART_CMD=systemctl --user restart nowify-kiosk
#   Environment=NOWIFY_RESTART_CMD=sudo systemctl restart nowify-kiosk
Environment=NOWIFY_RESTART_CMD=
ExecStart=/usr/bin/python3 /home/pi/Nowify/pi/nowify-watchdog.py
